# installer/systemd/actuator.service
[Unit]
Description=SLIME actuator-min (reference) — owns /run/slime/egress.sock
After=network.target
Wants=network.target

[Service]
Type=simple
User=actuator
Group=slime-actuator
RuntimeDirectory=slime
RuntimeDirectoryMode=0755
ExecStartPre=/bin/sh -c 'rm -f /run/slime/egress.sock'
ExecStartPre=/usr/local/bin/fireplank-guard-boot.sh
ExecStart=/usr/local/bin/actuator-min
Restart=no
NoNewPrivileges=yes
StandardOutput=append:/var/log/slime-actuator/events.log
StandardError=append:/var/log/slime-actuator/events.log

[Install]
WantedBy=multi-user.target
