# installer/systemd/slime-dashboard.service
[Unit]
Description=SLIME Dashboard (read-only observation) - noncanon enterprise surface
After=network.target actuator.service
Wants=actuator.service

[Service]
Type=simple
User=slime
Group=slime-actuator
EnvironmentFile=-/etc/slime/dashboard.env
ExecStart=/usr/bin/python3 /opt/slime/dashboard/server.py
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
ReadOnlyPaths=/var/log/slime-actuator /usr/lib/slime
# MiniMax audit finding MED-03: apply MemoryDenyWriteExecute uniformly
# across all three services (actuator and runner already have it via
# their fp4-hardening drop-ins). Python interpreters don't benefit from
# W^X the same way AOT-compiled code does, but consistent hardening
# makes the appliance surface easier to reason about and blocks any
# future JIT-enabled Python extension from introducing W+X mappings.
MemoryDenyWriteExecute=yes
# Copilot audit finding H-3: drop all capabilities. The dashboard is a
# stdlib HTTP server that only reads local files and makes outbound HTTP
# calls to a validated Ollama host; it needs no POSIX capabilities.
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
SystemCallFilter=@system-service

[Install]
WantedBy=multi-user.target
