# deploy/systemd/systemd/actuator.service
[Unit]
Description=SLIME actuator-min (reference) — owns /run/slime/egress.sock
After=network.target
Wants=network.target

[Service]
Type=simple
User=actuator
Group=slime-actuator
RuntimeDirectory=slime
RuntimeDirectoryMode=0750
# Kimi audit finding M-2: redundant rm removed — systemd's
# `RuntimeDirectory=slime` already provides a clean /run/slime on each
# service start (no `RuntimeDirectoryPreserve=`), so the explicit
# `rm -f /run/slime/egress.sock` added a TOCTOU window for no benefit.
ExecStartPre=/usr/local/bin/fireplank-guard-boot.sh
ExecStart=/usr/local/bin/actuator-min
Restart=no
NoNewPrivileges=yes
StandardOutput=append:/var/log/slime-actuator/events.log
StandardError=append:/var/log/slime-actuator/events.log

[Install]
WantedBy=multi-user.target
