# deploy/systemd/systemd/actuator.service.d/override.conf
[Service]
ExecStartPre=/bin/sh -lc 'rm -f /run/slime/egress.sock'
