# deploy/systemd/systemd/slime.service.d/override.conf
[Service]
WorkingDirectory=/opt/slime
ExecStartPre=/bin/sh -lc 'for i in $(seq 1 50); do test -S /run/slime/egress.sock && exit 0; sleep 0.1; done; exit 1'
StandardOutput=journal
StandardError=journal
